# recipes-mender/mender-standalone-daemon/files/mender-standalone-daemon.path
[Path]
DirectoryNotEmpty=@@MENDER/STANDALONE_DAEMON_DATA_DIR@@
MakeDirectory=True
DirectoryMode=1777

[Install]
WantedBy=multi-user.target
