# recipes-mender/mender-standalone-daemon/files/mender-standalone-daemon.service
[Unit]
Description=Mender Standalone Update Daemon

[Service]
Type=simple

RuntimeDirectory=mender-standalone-daemon

ExecStart=@@sbindir@@/mender-standalone-daemon.sh

[Install]
WantedBy=multi-user.target
